# sshuser_setup.service
[Unit]
Description=SSH user setup
After=sshd.service

[Service]
ExecStart=/bin/sshuser_setup
Type=oneshot

[Install]
WantedBy=default.target
